# imports/copilot/update-instruction-files/com.activitytracker.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.activitytracker.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/jacklicatamacbook/activity_tracker.py</string>
        <string>--daemon</string>
        <string>--poll</string>
        <string>2</string>
        <string>--idle</string>
        <string>300</string>
        <string>--heartbeat</string>
        <string>60</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/jacklicatamacbook</string>

    <key>StandardOutPath</key>
    <string>/Users/jacklicatamacbook/Library/Logs/ActivityTracker.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jacklicatamacbook/Library/Logs/ActivityTracker.err.log</string>
</dict>
</plist>
